# packages/base-files-lt70/files/etc/init.d/initscripts
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 teleofis

START=12
STOP=84

start()
{
	# modem init
	echo 495 > /sys/class/gpio/export
	echo high > /sys/class/gpio/gpio495/direction
	echo 480 > /sys/class/gpio/export
	echo low > /sys/class/gpio/gpio480/direction
	echo 493 > /sys/class/gpio/export
	echo in > /sys/class/gpio/gpio493/direction
	echo 494 > /sys/class/gpio/export
	echo in > /sys/class/gpio/gpio494/direction
	echo 496 > /sys/class/gpio/export
	echo low > /sys/class/gpio/gpio496/direction
	echo 497 > /sys/class/gpio/export
	echo low > /sys/class/gpio/gpio497/direction

	ln -s /dev/ttyS2 /dev/gps0
	ln -s /dev/ttyS1 /dev/com0
	ln -s /dev/ttyS0 /dev/com1
}
