# packages/base-files-common/files/etc/uci-defaults/30_ntpd
#!/bin/sh

[ ! -e "/usr/sbin/ntpd" && -e "/sbin/ntpd" ] || ln -s /sbin/ntpd /usr/sbin/ntpd
[ -e "/etc/init.d/sysntpd" ] && rm /etc/init.d/sysntpd 

exit 0
