# packages/base-files-rtux68/files/etc/uci-defaults/40_feeds
#!/bin/sh

mv /etc/opkg/distfeeds_cust.conf /etc/opkg/distfeeds.conf

exit 0
