# luci-app-simman/root/etc/uci-defaults/simman
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@simman[-1]
	add ucitrack simman
	set ucitrack.@simman[-1].init=simman
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
