# etc/init.d/telegram_bot
#!/bin/sh /etc/rc.common
# Copyright (c) 2011-2015 OpenWrt.org

START=99

start() {
	if [ ! -f "/var/run/telegram_bot.pid" ]; then
		echo "1" > /var/run/telegram_bot.pid
		/usr/lib/telegram-bot/telegram_bot >/dev/null 2>&1 &
	fi
	return 0
}

stop() {
	rm -f /var/run/telegram_bot.pid
	kill $(pgrep -f telegram-bot/telegram_bot) >/dev/null 2>&1
}

restart(){
	stop
	sleep 1
	start
}

boot() {
	(sleep 5 && start >/dev/null 2>&1) &
	return 0
}