# packages/net/openvpn-dns-hotplug/root/etc/hotplug.d/openvpn/90-dns
#!/bin/sh
# Get DNS, Domain Search via OpenVPN client hotplug script
# Written by Konstantine Shevlakov at <shevlakov@132lan.ru> 2025
# OpenWrt 21 or higher file location: /etc/hotplug.d/openvpn


# Get NameServers and Domain search
get_ns(){
	DNS=$(env |  awk '/dhcp-option DNS/{gsub("'\''", "");print $NF}')
	DOMAIN=$(env |  awk '/dhcp-option DOMAIN/{gsub("'\''", "");print $NF}')
}

# Setting up NameServers and Domain
set_ns(){
	if [ -z "$CFGIFACE" ]; then
		return 1
	fi
    
    # check interface in uci conf
    if ! ubus call network.interface.$CFGIFACE status >/dev/null 2>&1; then
        return 1
    fi
    
    # DNS processing
    if [ -n "$DNS" ]; then
        local dns_list=""
        for ns in $DNS; do
            [ -n "$ns" ] && dns_list="${dns_list:+$dns_list,}\"$ns\""
        done
        
        if [ -n "$dns_list" ]; then
            ubus call network.interface.$CFGIFACE set_dns "{\"dns\":[$dns_list]}" || {
                return 1
            }
        fi
    fi
    
    # domain search process
    if [ -n "$DOMAIN" ]; then
        local domain_list=""
        for dn in $DOMAIN; do
            [ -n "$dn" ] && domain_list="${domain_list:+$domain_list,}\"$dn\""
        done
        
        if [ -n "$domain_list" ]; then
            ubus call network.interface.$CFGIFACE set_search "{\"search\":[$domain_list]}" || {
                return 1
            }
        fi
    fi
}

# Delete NameServers and Domain
unset_ns(){
    [ -n "$CFGIFACE" ] || return 1
    ubus call network.interface.$CFGIFACE set_dns "{'dns':[]}" 2>/dev/null
    ubus call network.interface.$CFGIFACE set_search "{'search':[]}" 2>/dev/null
}


# Get interface name in uci 
UCI_conf=$(uci show openvpn | awk -F [\/=] '/config/{gsub("'\''", "");print $NF}')
for cfg in $UCI_conf; do
	if [ $config = $cfg ]; then
		CFGIFACE=$(uci show network | awk -F [.] -v device=$dev '$0 ~ device {print $2}')
		PROTO=$(ifstatus $CFGIFACE | jsonfilter -e '@["proto"]')
		get_ns
	fi
done

# Hotplug action
case $PROTO in 
	none)
		# Check peerdns option in uci interface section
		PEERDNS=$(uci -q get network.$CFGIFACE.dns)
		case $ACTION in
			up)
				if [ "x${PEERDNS}" != "x" ]; then
					exit 0
				fi
				if  [ "$DNS" ]; then
					unset_ns
					set_ns
					reload_config
					logger -t "`basename $0`" "OpenVPN: adding DNS: $DNS from $INSTANCE" 
				fi
			;;
			down)
				unset_ns
				reload_config
				logger -t "`basename $0`" "OpenVPN: remove NS attributes from $INSTANCE"
			;;
		esac
	;;
esac
