# luci/applications/luci-app-atinout/root/etc/uci-defaults/60_luci-app-atinout
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci set atinout.general=atinout
fi

if [ ! -f /etc/atcommands.user  ]; then
	echo "Demo command \"ATI\";ATI" > /etc/atcommands.user
fi

