# luci/applications/luci-app-atinout/root/etc/uci-defaults/65-luci-app-atinout
#!/bin/sh

if [ ! -f /etc/atcommands.user  ]; then
	echo "Demo command \"ATI\";ATI" > /etc/atcommands.user
fi

rm -rf /tmp/luci-*
