# luci/applications/luci-app-ttl/root/etc/hotplug.d/iface/90-ttl
sleep 5 && /etc/init.d/firewall restart &
