# luci/applications/luci-app-ttl/root/etc/hotplug.d/iface/90-ttl
. /lib/functions.sh

handle_config(){
	config_get iface "$1" iface
	[ "$iface" = "$INTERFACE" ] && {
		run_ttl=1
	}
}

config_load ttl
config_foreach handle_config ttl

case $run_ttl in
	1) /etc/init.d/ttl reload ;;
esac

exit 0

