# packages/base-files-rtux68/files/etc/hotplug.d/usb/12-com2
echo $DEVPATH | grep -q "/devices/soc0/80080000.ahb/80090000.usb/ci_hdrc.1/usb1/1-1/1-1.3/1-1.3:1.0" || exit 0
if [ "$PRODUCT" == "10c4/ea70/100" ]; then
	case "$ACTION" in
		add)
			ln -s /dev/$(ls /sys/$DEVPATH | grep tty) /dev/com2
			exit 0
		;;
		remove)
			rm /dev/com2
			exit 0
		;;
	esac
fi
